# ludusavi-backup.path
[Unit]
Description=Trigger Ludusavi Backup on Game State Change

[Path]
PathModified=/home/deck/.steam/steam/logs/controller.txt

[Install]
WantedBy=multi-user.target

